# applications/luci-qos/root/etc/uci-defaults/luci-qos
#!/bin/sh

/etc/init.d/qos enabled && {
	/etc/init.d/qos stop
	/etc/init.d/qos disable
}

rm -f /tmp/luci-indexcache
exit 0
