# applications/luci-samba/root/etc/uci-defaults/luci-samba
#!/bin/sh

/etc/init.d/samba enabled && {
	/etc/init.d/samba stop
	/etc/init.d/samba disable
}

rm -f /tmp/luci-indexcache
exit 0
